# reports/daily-order-report.timer
[Unit]
Description=Run Daily Order Report at HH:MM daily
Requires=daily-order-report.service

[Timer]
OnCalendar=*-*-* 08:00:00
Persistent=true

[Install]
WantedBy=timers.target